# templates/default/nerve.service.erb
[Unit]
Description=SmartStack Nerve
After=network.target

[Service]
Type=simple
Environment="GEM_HOME=<%= node.smartstack.gem_home %>" "RBENV_ROOT=/usr/local/rbenv"
WorkingDirectory=<%= node.nerve.install_dir %>
User=<%= node.smartstack.user %>
Group=<%= node.smartstack.user %>
Restart=on-failure
ExecStart=/usr/local/rbenv/shims/bundle exec <%= node.nerve.install_dir %>/bin/nerve --config <%= node.nerve.config_file %> 2>&1
